# bombinger.timer
[Timer]
OnCalendar=3:00
Persistent=false