# bombinger.timer
[Unit]
Requires=bombinger.service

[Timer]
OnCalendar=3:00
OnCalendar=14:00
Persistent=false

[Install]
WantedBy=timers.target
